# backend/com.gmu.dining.cron.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.gmu.dining.cron</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/Nagel/masondiningapp/backend/venv/bin/python</string>
        <string>/Users/Nagel/masondiningapp/backend/cron_job.py</string>
    </array>

    <!-- Run at 7:00 AM Eastern every day -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Eastern Time -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>TZ</key>
        <string>America/New_York</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/Nagel/masondiningapp/backend</string>

    <key>StandardOutPath</key>
    <string>/Users/Nagel/masondiningapp/backend/cron.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/Nagel/masondiningapp/backend/cron.log</string>

    <!-- Restart automatically if it crashes -->
    <key>KeepAlive</key>
    <false/>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
